# ops/systemd/publish-check.service.in
[Unit]
Description=Check that every due bessavagner.com post is pushed to origin/main

[Service]
Type=oneshot
ExecStart=@REPO@/ops/systemd/run-check.sh @REPO@ @NODE@
